# deploy/prim/services/prim-spawn.service
[Unit]
Description=prim.sh spawn.sh service
After=network.target

[Service]
Type=simple
User=prim
WorkingDirectory=/opt/prim
EnvironmentFile=/etc/prim/spawn.env
ExecStart=/home/prim/.bun/bin/bun run packages/spawn/src/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
