# deploy/prim/services/prim-treasury-refill.timer
[Unit]
Description=Daily treasury refill via Circle faucet

[Timer]
OnCalendar=*-*-* 06:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
