# deploy/prim/services/prim-faucet.service
[Unit]
Description=prim.sh faucet service
After=network.target

[Service]
Type=simple
User=prim
WorkingDirectory=/opt/prim
EnvironmentFile=/etc/prim/faucet.env
ExecStart=/home/prim/.bun/bin/bun run packages/faucet/src/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
