# deploy/prim/services/prim-search.service
[Unit]
Description=prim.sh search.sh service
After=network.target

[Service]
Type=simple
User=prim
WorkingDirectory=/opt/prim
EnvironmentFile=/etc/prim/search.env
ExecStart=/home/prim/.bun/bin/bun run packages/search/src/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
