# deploy/prim/services/prim-wallet.service
[Unit]
Description=prim.sh wallet.sh service
After=network.target

[Service]
Type=simple
User=prim
WorkingDirectory=/opt/prim
EnvironmentFile=/etc/prim/wallet.env
ExecStart=/home/prim/.bun/bin/bun run packages/wallet/src/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
