# deploy/prim/services/prim-feedback.service
[Unit]
Description=prim.sh feedback.sh service
After=network.target

[Service]
Type=simple
User=prim
WorkingDirectory=/opt/prim
EnvironmentFile=/etc/prim/feedback.env
ExecStart=/home/prim/.bun/bin/bun run packages/feedback/src/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
