# deploy/prim/services/prim-store.service
[Unit]
Description=prim.sh store.sh service
After=network.target

[Service]
Type=simple
User=prim
WorkingDirectory=/opt/prim
EnvironmentFile=/etc/prim/store.env
ExecStart=/home/prim/.bun/bin/bun run packages/store/src/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
