# deploy/prim/services/prim-gate.service
# THIS FILE IS GENERATED — DO NOT EDIT
# Source: packages/gate/prim.yaml
# Regenerate: pnpm gen:deploy

[Unit]
Description=prim.sh gate.sh service
After=network.target

[Service]
Type=simple
User=prim
WorkingDirectory=/opt/prim
EnvironmentFile=/etc/prim/gate.env
ExecStart=/home/prim/.bun/bin/bun run packages/gate/src/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
