# deploy/prim/services/prim-mem.service
[Unit]
Description=prim.sh mem.sh service
After=network.target

[Service]
Type=simple
User=prim
WorkingDirectory=/opt/prim
EnvironmentFile=/etc/prim/mem.env
ExecStart=/home/prim/.bun/bin/bun run packages/mem/src/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
