# deploy/prim/services/prim-site.service
[Unit]
Description=prim.sh site server
After=network.target

[Service]
Type=simple
User=prim
WorkingDirectory=/opt/prim
ExecStart=/home/prim/.bun/bin/bun run site/serve.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
